# tools/local_celerybeat.conf
; ================================
;  celery beat supervisor example
; ================================

[program:celerybeat]
; Set full path to celery program if using virtualenv
command=celery beat -A databoard:celery -s DATABOARD_PATH/code/databoard/tools/celery_info/celerybeat-schedule --loglevel=INFO --pidfile="DATABOARD_PATH/code/databoard/tools/celery_info/celerybeat.pid"

; remove the -A myapp argument if you are not using an app instance

directory=DATABOARD_PATH/code/databoard
user=ubuntu
numprocs=1
stdout_logfile=DATABOARD_PATH/code/databoard/tools/celery_info/beat.log
stderr_logfile=DATABOARD_PATH/code/databoard/tools/celery_info/beat.log
autostart=true
autorestart=true
startsecs=10

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999

